# systemd/agnetic-dashboard.service
[Unit]
Description=Starship OS - Web Dashboard (:8788)
After=network.target agnetic-nats.service
Requires=agnetic-nats.service
Documentation=https://github.com/andromi-hash/starship-os

[Service]
Type=simple
ExecStart=/opt/starship/venv/bin/python3 /opt/starship/lib/starship/dashboard/server.py
Restart=always
RestartSec=5
User=agnetic
Group=agnetic
EnvironmentFile=-/etc/starship/nats.env
Environment=PYTHONPATH=/opt/starship/lib/starship:/opt/starship/lib/starship/agents
Environment=STARSHIP_ROOT=/opt/starship
Environment=AGNETIC_ROOT=/opt/starship
Environment=DASHBOARD_PORT=8788
Environment=AGNETIC_DASHBOARD_PORT=8788
WorkingDirectory=/opt/starship/lib/starship
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log/starship /var/log/agnetic /tmp
ProtectHome=true
KillMode=mixed
TimeoutStopSec=10
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
