# systemd/agnetic-agent@.service
[Unit]
Description=Starship OS - Agent Daemon (%i)
After=network.target agnetic-nats.service
Requires=agnetic-nats.service
Documentation=https://github.com/andromi-hash/starship-os

[Service]
Type=simple
ExecStart=/opt/starship/venv/bin/python3 /opt/starship/lib/starship/agents/agent_daemon.py %i
Restart=always
RestartSec=10
User=agnetic
Group=agnetic
EnvironmentFile=-/etc/starship/nats.env
EnvironmentFile=-/etc/starship/fleet.env
Environment=PYTHONPATH=/opt/starship/lib/starship:/opt/starship/lib/starship/agents
Environment=STARSHIP_ROOT=/opt/starship
Environment=AGNETIC_ROOT=/opt/starship
Environment=OLLAMA_URL=http://127.0.0.1:11434
Environment=NATS_URL=nats://127.0.0.1:4222
WorkingDirectory=/opt/starship/lib/starship
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log/starship /var/log/agnetic /tmp
ProtectHome=true
KillMode=mixed
TimeoutStopSec=15
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
