# deploy/agnetic-status-bridge.service
[Unit]
Description=Starship OS Status Bridge — NATS telemetry to JSON file
After=agnetic-nats.service
Requires=agnetic-nats.service
PartOf=agnetic-dashboard.target

[Service]
Type=simple
User=agnetic
Environment=NATS_URL=nats://127.0.0.1:4222
Environment=STARSHIP_ROOT=/opt/starship
Environment=AGNETIC_ROOT=/opt/starship
ExecStart=/opt/starship/venv/bin/python3 /opt/starship/lib/starship/tray/agnetic-status.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=agnetic-dashboard.target
