# deploy/agnetic-dashboard-web.service
[Unit]
Description=Starship OS Web Dashboard — HTTP server on port 8788
After=agnetic-status-bridge.service
Requires=agnetic-status-bridge.service
PartOf=agnetic-dashboard.target

[Service]
Type=simple
User=agnetic
Environment=DASHBOARD_PORT=8788
Environment=AGNETIC_DASHBOARD_PORT=8788
Environment=STARSHIP_ROOT=/opt/starship
Environment=AGNETIC_ROOT=/opt/starship
ExecStart=/opt/starship/venv/bin/python3 /opt/starship/lib/starship/dashboard/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=agnetic-dashboard.target
