# config/osquery/starshipd.conf
{
  "options": {
    "host_identifier": "uuid",
    "schedule_splay_percent": 10,
    "pidfile": "/var/run/osqueryd.pid",
    "database_path": "/var/lib/starship/osquery/db",
    "logger_path": "/var/log/starship/osquery",
    "logger_mode": "json",
    "logger_event_type": true,
    "disable_logging": false,
    "disable_events": false,
    "event_pubsub": true,
    "events_expiry": 3600,
    "enable_syslog": false,
    "aws_firehose_period": 10,
    "aws_firehose_stream": "",
    "allow_unsafe": false,
    "audit_allow_config": true,
    "audit_allow_sockets": true,
    "distributed_interval": 60,
    "distributed_tls_read_endpoint": "",
    "distributed_tls_write_endpoint": "",
    "pack_refresh_interval": 3600,
    "schedule_default_interval": 3600,
    "schedule_timeout": 60,
    "utc": true,
    "verbose": false,
    "worker_threads": 2,
    "watchdog_level": 1,
    "enable_mount_events": true,
    "enable_system_events": true,
    "enable_filesystem_events": true,
    "events_max": 50000
  },
  "schedule": {
    "system_info": {
      "query": "SELECT hostname, cpu_brand, cpu_logical_cores, physical_memory, hardware_vendor, hardware_model FROM system_info;",
      "interval": 3600,
      "description": "System hardware and OS info"
    },
    "os_version": {
      "query": "SELECT name, version, major, minor, patch, build, platform FROM os_version;",
      "interval": 86400,
      "description": "OS version details"
    },
    "process_list": {
      "query": "SELECT pid, name, path, cmdline, cpu_percent, resident_size, total_size, user_time, system_time, start_time, parent FROM processes ORDER BY resident_size DESC;",
      "interval": 300,
      "description": "Top processes by memory usage"
    },
    "listening_ports": {
      "query": "SELECT pid, port, protocol, address, fd, socket, path FROM listening_ports WHERE port > 0;",
      "interval": 300,
      "description": "Listening network ports"
    },
    "process_open_sockets": {
      "query": "SELECT pid, fd, socket, protocol, local_address, local_port, remote_address, remote_port, state FROM process_open_sockets WHERE remote_address != '';",
      "interval": 300,
      "description": "Active network connections"
    },
    "logged_in_users": {
      "query": "SELECT user, type, tty, host, time, pid FROM logged_in_users;",
      "interval": 120,
      "description": "Currently logged-in users"
    },
    "last_boot": {
      "query": "SELECT * FROM last_boot;",
      "interval": 3600,
      "description": "System uptime and last boot"
    },
    "mounts": {
      "query": "SELECT device, path, type, blocks_size, blocks, blocks_available, inodes, inodes_free FROM mounts WHERE path != '';",
      "interval": 600,
      "description": "Filesystem mount points"
    },
    "disk_encryption": {
      "query": "SELECT name, type, encrypted, encryption_type, user_uuid FROM disk_encryption;",
      "interval": 3600,
      "description": "Disk encryption status"
    },
    "kernel_modules": {
      "query": "SELECT name, size, used_by, status, address FROM kernel_modules;",
      "interval": 600,
      "description": "Loaded kernel modules"
    },
    "docker_containers": {
      "query": "SELECT id, name, image, state, status, pid, started_at, command, ports FROM docker_containers;",
      "interval": 300,
      "description": "Running Docker containers"
    },
    "crontab": {
      "query": "SELECT event, command, minute, hour, day_of_month, month, day_of_week FROM crontab;",
      "interval": 600,
      "description": "Scheduled cron jobs"
    }
  },
  "packs": [
    "starship_security",
    "starship_compliance"
  ],
  "file_paths": {
    "etc_passwd": [
      "/etc/passwd"
    ],
    "etc_shadow": [
      "/etc/shadow"
    ],
    "etc_sudoers": [
      "/etc/sudoers",
      "/etc/sudoers.d/"
    ],
    "ssh_config": [
      "/etc/ssh/ssh_config",
      "/etc/ssh/sshd_config",
      "/root/.ssh/",
      "/home/%%/.ssh/"
    ],
    "system_binaries": [
      "/usr/bin/",
      "/usr/sbin/",
      "/bin/",
      "/sbin/"
    ],
    "log_files": [
      "/var/log/"
    ]
  },
  "file_accesses": {
    "restricted": [
      "/etc/shadow",
      "/etc/sudoers",
      "/etc/ssh/ssh_host_*"
    ]
  },
  "views": {
    "full_system_info": "select * from system_info limit 1;"
  }
}
